# data/sd/mirrors-live-net.service
[Unit]
Description=LiveMedia Pacman mirrors ranking script
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/pacman-mirrors -f 5

[Install]
WantedBy=multi-user.target
